# packaging/systemd/podlazd.service
[Unit]
Description=podlaz daemon
After=network.target
